# sauna-gateway/rpi-activity-led-gpio.service
[Unit]
Description=Change RPi activity led to gpio mode

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo none > /sys/class/leds/ACT/trigger'
ExecStart=/bin/sh -c 'chown root.gpio /sys/class/leds/ACT/brightness'
ExecStart=/bin/sh -c 'chmod g+rw /sys/class/leds/ACT/brightness'
User=root

[Install]
WantedBy=multi-user.target
Alias=rpi-activity-led-gpio.service
